# roles/plex-client/files/x11.conf
start on filesystem
stop on runlevel [06]

exec X -dpms
